# ansible/roles/app/templates/app.service.j2
[Unit]
Description=Smoke app de validacion para la Fase 4
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ phase4_app_user }}
Group={{ phase4_app_group }}
WorkingDirectory={{ phase4_app_dir }}
EnvironmentFile={{ phase4_app_dir }}/.env
ExecStart=/usr/bin/node {{ phase4_app_dir }}/server.js
StandardOutput=append:{{ phase4_app_log_file }}
StandardError=append:{{ phase4_app_log_file }}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
