# ansible/roles/backend/templates/cotizacion-app.service.j2
[Unit]
Description=Cotizacion Backend Application
After=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/opt/phase4-app/current
EnvironmentFile=/opt/phase4-app/current/.env
ExecStart=/usr/bin/node start.js
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cotizacion-app

# Temporarily relaxed security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/phase4-app

[Install]
WantedBy=multi-user.target
