# ansible/roles/backend/templates/cotizacion-app.service.j2
[Unit]
Description=Cotizacion Backend Application
After=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
# WorkingDirectory usa la ruta real del release, no el symlink,
# para que Node encuentre dist/server.js con require("./dist/server.js")
WorkingDirectory={{ backend_release_path }}
EnvironmentFile={{ backend_release_path }}/.env
ExecStart=/usr/bin/node {{ backend_release_path }}/start.js
Restart=always
RestartSec=10

# Logging a journal — accesible con: journalctl -u cotizacion-app -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cotizacion-app

# Seguridad mínima
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths={{ backend_release_path }} /opt/phase4-app /var/log

[Install]
WantedBy=multi-user.target
